# debian/system76-power.postinst
#!/bin/sh

set -e

systemctl enable com.system76.PowerDaemon

case "$1" in
    configure)
        rm -f /etc/modules-load.d/system76-power.conf
        ;;

    *)
        ;;
esac

#DEBHELPER#

exit 0
